# crontab
0 23 * * * root /cron_task.sh /pg_backup/pg_backup_rotated.sh
